# etc/graphite/supervisor.statsd.conf
[program:statsd]
command = /usr/bin/node stats.js config.js
directory = /opt/statsd/
user = tom
autostart = true
autorestart = true
redirect_stderr = true