# vagrant/etc/graphite/supervisor.carbon.conf
[program:carbon]
command = /opt/graphite/bin/carbon-cache.py start --debug
user = tom
autostart = true
autorestart = true
redirect_stderr = true